# IAM.socket
[Unit]
Description=IAM socket

[Socket]
ListenStream=/run/IAM.sock

[Install]
WantedBy=sockets.target
